# .config/systemd/user/home-tonyampomah-Mount-Shared.automount
[Unit]
  Description=Automount Downloads

[Automount]
  Where=/home/tonyampomah/Mount/Shared

[Install]
  WantedBy=multi-user.target
