# .config/systemd/user/home-tonyampomah-Mount-Iso.automount
[Unit]
  Description=Automount Iso

[Automount]
  Where=/home/tonyampomah/Mount/Iso

[Install]
  WantedBy=multi-user.target
