# .config/systemd/user/home-kwamedat-Mount-Pictures.automount
[Unit]
  Description=Automount Pictures

[Automount]
  Where=/home/kwamedat/Mount/Pictures

[Install]
  WantedBy=multi-user.target
