# .config/systemd/user/home-kwamedat-Mount-Iso.automount
[Unit]
  Description=Automount Iso

[Automount]
  Where=/home/kwamedat/Mount/Iso

[Install]
  WantedBy=multi-user.target
