# .config/systemd/user/home-tonyampomah-Mount-Pictures.automount
[Unit]
  Description=Automount Pictures

[Automount]
  Where=/home/tonyampomah/Mount/Pictures

[Install]
  WantedBy=multi-user.target
