# .config/systemd/user/home-tonyampomah-Mount-Videos.automount
[Unit]
  Description=Automount Videos

[Automount]
  Where=/home/tonyampomah/Mount/Videos

[Install]
  WantedBy=multi-user.target
