# .config/systemd/user/home-kwamedat-Mount-Videos.automount
[Unit]
  Description=Automount Videos

[Automount]
  Where=/home/kwamedat/Mount/Videos

[Install]
  WantedBy=multi-user.target
