# config/systemd/user/home-kwamedat-Mount-Shared.automount
[Unit]
  Description=Automount Downloads

[Automount]
  Where=/home/kwamedat/Mount/Shared

[Install]
  WantedBy=multi-user.target
